# templates/launchd/com.user.codingtracker.stats.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.codingtracker.stats</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>__SCRIPT_DIR__/scripts/push_stats.sh</string>
    </array>
    <key>WorkingDirectory</key>
    <string>__SCRIPT_DIR__</string>
    <key>WatchPaths</key>
    <array>
        <string>__HOME__/.coding_tracker.db</string>
        <string>__HOME__/.coding_tracker.db-wal</string>
        <string>__HOME__/.coding_tracker.db-shm</string>
    </array>
    <key>ThrottleInterval</key>
    <integer>300</integer>
    <key>RunAtLoad</key>
    <true/>
    <key>StartInterval</key>
    <integer>3600</integer>
    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Hour</key>
            <integer>13</integer>
            <key>Minute</key>
            <integer>30</integer>
        </dict>
        <dict>
            <key>Hour</key>
            <integer>20</integer>
            <key>Minute</key>
            <integer>30</integer>
        </dict>
    </array>
    <key>StandardOutPath</key>
    <string>__HOME__/.coding_tracker_stats_stdout.log</string>
    <key>StandardErrorPath</key>
    <string>__HOME__/.coding_tracker_stats_stderr.log</string>
</dict>
</plist>
